# toolbox/auto-update/.config/systemd/user/update-container@.timer
[Unit]
Description=Update user Toolbox daily (%i)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target